# extras/launchagents/com.typeshield.autorestart-smart.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.typeshield.autorestart-smart</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__HOME__/Library/Scripts/typeshield-auto-restart.sh</string>
  </array>

  <key>StartInterval</key>
  <integer>5</integer>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/ts-smart.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/ts-smart.err</string>
</dict>
</plist>
